# android.service
[Unit]
Description=Android in a chroot

[Service]
ExecStart=/init
RootDirectory=/android
TimeoutSec=10s

[Install]
WantedBy=graphical.target
